# scraper/scrape-schedule.cron
# Run annual report scraper once per week on early Saturday morning.
 5 4 * * 6 /home/rstudio/repel-infrastructure/scraper/scrape-annual-reports.R
# Run run outbreak report scraper daily
 5 0 * * * /home/rstudio/repel-infrastructure/scraper/scrape-outbreak-reports.R
# We need a blank line here for it to be a valid cron file
